# mac-mini/LaunchAgents/com.mizuta.claude-rc.kanken.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mizuta.claude-rc.kanken</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="$HOME/.local/bin:$PATH"; claude remote-control --spawn worktree --sandbox</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/mizuta/work/products/kanken</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/mizuta/Library/Logs/claude-rc/kanken.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mizuta/Library/Logs/claude-rc/kanken.err.log</string>
</dict>
</plist>
