# mac-mini/LaunchAgents/com.mizuta.colima.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mizuta.colima</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="/opt/homebrew/bin:$PATH"; colima start --cpus 8 --memory 16</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/mizuta/Library/Logs/colima/colima.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mizuta/Library/Logs/colima/colima.err.log</string>
</dict>
</plist>
